# cardloop.service.template
# systemd unit template for Cardloop.
#
# Easiest: run `make service` (renders this template for the current user,
# installs it, then enables & starts it). Manual install:
#   1. Substitute the username:
#        sed 's/__USER__/YOUR_USERNAME/g' cardloop.service.template > /tmp/cardloop.service
#   2. sudo cp /tmp/cardloop.service /etc/systemd/system/cardloop.service
#   3. sudo systemctl daemon-reload && sudo systemctl enable --now cardloop

[Unit]
Description=Cardloop — Claude Code IDE (web cockpit + kanban auto-run)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=%h/cardloop
Environment=HOME=%h
Environment=PATH=%h/.local/bin:%h/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=%h/cardloop/.env
# ANTHROPIC_API_KEY is intentionally unset — we use Claude subscription OAuth (~/.claude/.credentials.json)
ExecStart=%h/cardloop/venv/bin/python -u %h/cardloop/bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
