# cardloop.service.template
# systemd unit template for Cardloop.
#
# Easiest: run `make service` (renders this template for the current user,
# installs it, then enables & starts it). Manual install:
#   1. Substitute the username AND home dir (systemd's %h is unreliable for
#      system units with User= — it resolves to /root, not the user's home):
#        sed -e "s/__USER__/$(whoami)/g" -e "s#__HOME__#$HOME#g" \
#            cardloop.service.template > /tmp/cardloop.service
#   2. sudo cp /tmp/cardloop.service /etc/systemd/system/cardloop.service
#   3. sudo systemctl daemon-reload && sudo systemctl enable --now cardloop

[Unit]
Description=Cardloop — Claude Code IDE (web cockpit + kanban auto-run)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__HOME__/cardloop
Environment=HOME=__HOME__
Environment=PATH=__HOME__/.local/bin:__HOME__/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=__HOME__/cardloop/.env
# ANTHROPIC_API_KEY is intentionally unset — we use Claude subscription OAuth (~/.claude/.credentials.json)
ExecStart=__HOME__/cardloop/venv/bin/python -u __HOME__/cardloop/bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
